# bin/10-md407.rules
ATTR{product}=="FT230X Basic UART", RUN+="/usr/local/bin/setMD407baud"